# rootfs/etc/s6-overlay/s6-rc.d/init-services/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

ONESHOT=$(bashio::config 'Oneshot')
                          
if [ "$ONESHOT" != "true" ]; then
  for T in filter nat mangle raw security; do
    iptables -t $T -F || true
  done
fi
true
